# systemd/savalet.service
[Unit]
Description=Savalet Command Execution Daemon
Documentation=https://github.com/zinrai/savalet
After=network.target

[Service]
Type=notify
ExecStart=/usr/local/bin/savalet daemon --config /etc/savalet/daemon.yaml
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
Restart=on-failure
RestartSec=5s

# Security settings
User=savalet
Group=savalet
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/run
ReadOnlyPaths=/etc/savalet

# Resource limits
LimitNOFILE=4096
TasksMax=128

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=savalet

[Install]
WantedBy=multi-user.target
